# systemd_files/arch_updater.service
[Unit]
Description="ArchLinux updater"

[Service]
ExecStart=/usr/local/bin/arch_updater.sh
