# systemd_files/etc/systemd/system/arch_updater.service
[Unit]
Description="Hello World script"

[Service]
ExecStart=/usr/local/bin/helloworld.sh
